# scripts/stockyardd.service
[Unit]
Description=Stockyard VM Daemon
Documentation=https://github.com/obra/stockyard
After=network.target zfs.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/stockyardd
Restart=on-failure
RestartSec=5

# Runtime directory for socket
RuntimeDirectory=stockyard
RuntimeDirectoryMode=0755

# State directory for VMs and data
StateDirectory=stockyard
StateDirectoryMode=0755

# Stockyard needs broad system access for:
# - Firecracker VMs (requires /dev access, process management)
# - ZFS datasets (requires zfs commands)
# - Network bridges and TAP interfaces
# - DHCP server (dnsmasq binds to port 67)
# Minimal sandboxing is applied

[Install]
WantedBy=multi-user.target
