# vm-image/init/stockyard-init.service
[Unit]
Description=Stockyard VM Initialization
# Disable default dependencies so we can start as early as possible.
# Kernel configures static IP via boot args before systemd even starts.
DefaultDependencies=no
# Only wait for basic filesystem and network interface to exist
After=local-fs.target systemd-networkd.service
# Block sysinit.target until we're done - ensures we run early
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stockyard-init.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
# Pull in via sysinit.target for early execution
WantedBy=sysinit.target
