# vm-image/init/stockyard-shell.service
[Unit]
Description=Stockyard Shell Service (vsock terminal access)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/stockyard-shell
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal

# Graceful shutdown
TimeoutStopSec=5
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
