# vm-image/init/stockyard-init.initd
#!/sbin/openrc-run
# OpenRC init script for stockyard VM initialization
# Runs once at boot to configure the VM from MMDS metadata

description="Stockyard VM Initialization"

depend() {
    # Run before sshd so SSH keys are in place
    before sshd
    before stockyard-shell
    before llm-proxy
}

start() {
    # Guard against running multiple times
    if [ -f /run/stockyard-init.done ]; then
        einfo "Stockyard init already completed, skipping"
        return 0
    fi

    ebegin "Running Stockyard VM initialization"
    /usr/local/bin/stockyard-init.sh
    local ret=$?
    touch /run/stockyard-init.done
    eend $ret
}

stop() {
    # Nothing to stop — oneshot service
    return 0
}
